# recipes-support/openhr20/files/openhr20-backup.timer
[Unit]
Description=Backup OpenHR20 Databases once an hour

[Timer]
OnBootSec=1hour
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target