# recipes-support/openhr20/files/openhr20-plot-7days.timer
[Unit]
Description=Plot 3 day graph every 15min

[Timer]
OnBootSec=1minute
OnUnitActiveSec=1hour

[Install]
WantedBy=timers.target