# recipes-support/openhr20/files/openhr20-setttyspeed.service
[Unit]
Description=OpenHR20 Set TTY-Speed
Wants=dev-ttyUSB0.device
After=dev-ttyUSB0.device

[Service]
Type=oneshot
ExecStart=/bin/stty -F /dev/ttyUSB0 38400

[Install]
WantedBy=multi-user.target
